# etc/supervisor/conf.d/http_fs.conf
[program:http_fs]
directory=/opt/http_fs
command=python service.py --config=/etc/http_fs.conf